# jupyter.conf
[program:jupyter]
command=jupyter notebook --allow-root
user=root

